# allegro_hand_core_template/build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/nxr/allegro_ws/allegro_hand_ros/allegro_hand_controllers/allegro_hand_core_template/build/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
